# pmxt_relay/systemd/pmxt-disable-wbt.service
[Unit]
Description=Disable block writeback throttling for relay disks
After=local-fs.target
Before=clickhouse-server.service pmxt-relay-api.service pmxt-relay-worker.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for dev in xvda xvdb xvdc; do if [ -w /sys/block/$dev/queue/wbt_lat_usec ]; then echo 0 > /sys/block/$dev/queue/wbt_lat_usec; fi; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
